# com.booker.evolve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.booker.evolve</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/tanmay/IdeaProjects/kingpinXD/booker-evolve/scripts/evolve-cron.sh</string>
	</array>

	<key>StartInterval</key>
	<integer>14400</integer>

	<key>Nice</key>
	<integer>10</integer>

	<key>RunAtLoad</key>
	<false/>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/tanmay/.local/bin:/opt/homebrew/opt/go/libexec/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>/Users/tanmay</string>
		<key>GOPATH</key>
		<string>/Users/tanmay/go</string>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/tanmay/IdeaProjects/kingpinXD/booker-evolve/logs/launchd-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/tanmay/IdeaProjects/kingpinXD/booker-evolve/logs/launchd-stderr.log</string>
</dict>
</plist>
